# deploy/systemd/footballvision-api-enhanced.service
[Unit]
Description=FootballVision Pro Enhanced API (Resilient)
After=network-online.target
Wants=network-online.target
# Ensure we start after critical services
After=systemd-resolved.service
After=NetworkManager.service
# Avoid start-limit lockouts during frequent deploy/restart operations.
StartLimitIntervalSec=120
StartLimitBurst=30

[Service]
Type=simple
User=mislav
Group=mislav
WorkingDirectory=/home/mislav/footballvision-pro

# Environment setup
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/mislav/footballvision-pro/src"
Environment="PYTHONUNBUFFERED=1"

# Main service
ExecStart=/usr/bin/python3 /home/mislav/footballvision-pro/src/platform/simple_api_v3.py

# Pre-start actions
ExecStartPre=/usr/bin/mkdir -p /var/log/footballvision/api /var/log/footballvision/system /var/log/footballvision/crashes
# Optional hardware tuning: never block API startup if unavailable/permission-restricted.
ExecStartPre=-/bin/bash -lc 'command -v nvpmodel >/dev/null 2>&1 && nvpmodel -m 1'
ExecStartPre=-/bin/bash -lc 'for g in /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor; do [ -w "$g" ] && echo performance > "$g" || true; done'

# Post-start recovery check
# Recovery check must not fail service startup.
ExecStartPost=-/bin/bash /home/mislav/footballvision-pro/scripts/power_loss_recovery.sh

# Restart policy - aggressive but with backoff
Restart=always
RestartSec=5

# Kill settings - graceful shutdown for recordings
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Resource limits
LimitNOFILE=65536
Nice=-5

# Logging
StandardOutput=append:/var/log/footballvision/api/service.log
StandardError=append:/var/log/footballvision/api/error.log

# Security (minimal restrictions for camera access)
PrivateTmp=false
ProtectHome=false
ProtectSystem=false

[Install]
WantedBy=multi-user.target
